# deploy/memi-lx-webhook.service
[Unit]
Description=Memi Lisboa deploy webhook
After=network.target

[Service]
User=root
WorkingDirectory=/opt/memi-lx
EnvironmentFile=/etc/memi-lx-webhook.env
ExecStart=/opt/memi-lx/.venv/bin/python deploy/webhook.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
